# overlay/etc/skel/.config/autostart/nmapplet.desktop
[Desktop Entry]

Type=Application

Exec=nm-applet