# overlay/etc/skel/.config/autostart/parcellite.desktop
[Desktop Entry]

Type=Application

Exec=parcellite