# overlay/etc/skel/.config/autostart/release.desktop
[Desktop Entry]

Type=Application

Exec=firefox "file:///usr/share/doc/archbang/Release Notes.html"