# overlay/etc/skel/.config/autostart/volume.desktop
[Desktop Entry]

Type=Application

Exec=volwheel